# etc/service/dnsmasq/run
#!/bin/sh

# source environment variables
. /etc/default-environment
. /etc/environment

# change working directory
cd $(dirname $0)

# overwrite distribution configuration
echo -n >/etc/dnsmasq.conf

# create data directory
mkdir -p ${DNSMASQ_DATADIR}

# create hosts directory
mkdir -p ${DNSMASQ_DATADIR}/hosts

# generate resolv.conf file
# if NetworkManager is mapped, use its nameservers, or use the docker ones instead

resolvconf_file=${DNSMASQ_DATADIR}/resolv.conf.generated

if [ -f /run/NetworkManager/resolv.conf ]; then
	grep "^nameserver\|^search" /run/NetworkManager/resolv.conf >${resolvconf_file}
else
	grep "^nameserver\|^search" /etc/resolv.conf >${resolvconf_file}
fi

# fix /etc/resolv.conf which have been generated by docker
printf '1,$d\na\nnameserver 127.0.0.1\n.\nw\nq\n' | ed /etc/resolv.conf

exec dnsmasq \
	--keep-in-foreground \
	--conf-dir=${DNSMASQ_CONFDIR} \
	--pid-file=${DNSMASQ_DATADIR}/dnsmasq.pid \
	2>&1
