# etc/service/dnsmasq/run
#!/bin/sh

# source environment variables
. /etc/default-environment
. /etc/environment

# overwrite distribution configuration
echo -n >/etc/dnsmasq.conf

# create data directory
mkdir -p ${DNSMASQ_DATADIR}

# create hosts directory
mkdir -p ${DNSMASQ_DATADIR}/hosts

# create dynamic config directory
mkdir -p ${DNSMASQ_DATADIR}/conf

# fix /etc/resolv.conf which have been generated by docker
printf '1,$d\na\nnameserver 127.0.0.1\n.\nw\nq\n' | ed /etc/resolv.conf

# start dnsmasq
exec dnsmasq \
	--keep-in-foreground \
	--conf-dir=${DNSMASQ_CONFDIR} \
	--pid-file=${DNSMASQ_DATADIR}/dnsmasq.pid \
	2>&1
